# launchd/com.austenstone.thinking-phrases.rss.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.austenstone.thinking-phrases.rss</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>/Users/austenstone/source/vscode-copilot-tips/scripts/run-rss-update.zsh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>/Users/austenstone/source/vscode-copilot-tips/launchd/rss-update.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/austenstone/source/vscode-copilot-tips/launchd/rss-update.error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/austenstone/source/vscode-copilot-tips</string>
  </dict>
</plist>